# backend-php/websocket/gende-ws.service
[Unit]
Description=Gende WebSocket Server
After=network.target mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/backend-php/websocket
ExecStart=/usr/bin/php server.php
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=WS_PORT=8090
EnvironmentFile=-/var/www/backend-php/config/.env

[Install]
WantedBy=multi-user.target
